# scripts/com.openclaw.omniclaw-smoke.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.omniclaw-smoke</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/markshteyn/conductor/workspaces/omniclaw-v1/boston/scripts/run-smoke.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/Users/markshteyn/.openclaw/smoke-logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/markshteyn/.openclaw/smoke-logs/launchd-stderr.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/markshteyn</string>
    </dict>
</dict>
</plist>
